# preston-cli/src/main/deb/lib/systemd/system/preston-registration.timer
[Unit]
Description=Run Preston Registration

[Timer]
OnCalendar=weekly
RandomizedDelaySec=86400

[Install]
WantedBy=timers.target
