# src/main/deb/etc/systemd/system/preston.timer
[Unit]
Description=Run Preston

[Timer]
OnCalendar=weekly
RandomizedDelaySec=86400
